# rootd.service
[Unit]
Description=Root Daemon for GPT Control (rootd)
After=network.target

[Service]
WorkingDirectory=/home/admin/gptadmin
# Запуск из вашего .venv
ExecStart=/home/admin/gptadmin/.venv/bin/python /home/admin/gptadmin/rootd.py

# Переменные окружения: агентов регистрируются у Hub по HTTPS
Environment=ROOTD_TOKEN=srv_secret
Environment=ROOTD_URL=http://$(hostname):25900
Environment=HUB_URL=https://gptadmin.bezrabotnyi.com/heartbeat
Environment=EXEC_TIMEOUT=300
Environment=LOG_LIMIT_B=8192
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
