# deploy/systemd/gptadmin-cleanup.service
[Unit]
Description=GPTAdmin cleanup old outputs and logs
After=network.target

[Service]
Type=oneshot
User=admin
Group=admin
WorkingDirectory=/home/admin/gptadmin
ExecStart=/home/admin/gptadmin/scripts/cleanup.sh
StandardOutput=journal
StandardError=journal
