# deploy/systemd/gptadmin-hub-watchdog.service
[Unit]
Description=GPTAdmin Hub external watchdog
Documentation=man:systemd.service(5)
Wants=network-online.target
After=network-online.target hub_proxy.service

[Service]
Type=oneshot
WorkingDirectory=/home/admin/gptadmin
ExecStart=/home/admin/gptadmin/.venv/bin/python /home/admin/gptadmin/hub_watchdog.py --check-once
