# deploy/systemd/gptadmin-cleanup.service
[Unit]
Description=GPTAdmin cleanup old outputs and logs
After=network.target

[Service]
Type=oneshot
User=roomhacker
Group=roomhacker
WorkingDirectory=/home/roomhacker/gptadmin
ExecStart=/home/roomhacker/gptadmin/scripts/cleanup.sh
StandardOutput=journal
StandardError=journal
