# deploy/systemd/gptadmin-hub-watchdog.timer
[Unit]
Description=Run GPTAdmin Hub watchdog every 30 seconds

[Timer]
OnBootSec=45s
OnUnitActiveSec=30s
AccuracySec=5s
Persistent=true
Unit=gptadmin-hub-watchdog.service

[Install]
WantedBy=timers.target
