# deploy/systemd/rootd.service
[Unit]
Description=Root Daemon for GPT Control (legacy Python rootd, deprecated)
After=network.target

[Service]
# DEPRECATED: this template starts legacy Python rootd.
# Use go-shellmcp/rootd-go-canary for new deployments.
WorkingDirectory=/home/admin/gptadmin
# Запуск из вашего .venv
ExecStart=/home/admin/gptadmin/.venv/bin/python /home/admin/gptadmin/client/rootd.py

# Переменные окружения: агентов регистрируются у Hub по HTTPS
Environment=ROOTD_TOKEN=srv_secret
Environment=ROOTD_URL=http://203.0.113.10:25900
Environment=HUB_URL=https://gptadmin.bezrabotnyi.com/heartbeat
Environment=EXEC_TIMEOUT=300
Environment=LOG_LIMIT_B=10485760
Environment=ROOTD_TRANSPORT=webhook
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
