# deploy/systemd/shellmcp.service
[Unit]
Description=GPTAdmin ShellMCP Go transport
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/admin/gptadmin/go-shellmcp
EnvironmentFile=/etc/gptadmin/go-shellmcp-primary.env
ExecStart=/opt/gptadmin/bin/rootd-go-canary
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
